# feishu_sync_service_backup_20260221_1509/openclaw-feishu-sync.service
[Unit]
Description=Feishu Local Sync Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/n150/openclaw_workspace/feishu_sync_service
ExecStart=/usr/bin/python3 /home/n150/openclaw_workspace/feishu_sync_service/app.py daemon
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
